# contrib/containerized-pbench/dockerfiles/agent/pbench.service
[Unit]
Description=pbench daemon
After=sshd.service

[Service]
ExecStart=/root/mount.sh
StandardOutput=null

[Install]
WantedBy=multi-user.target
