# rpm/benchmarks/uperf/uperf.service
[Unit]
Description=uperf A network performance tool

[Service]
ExecStart=/usr/local/bin/uperf -s

[Install]
WantedBy=multi-user.target
