# bgtasks/pbench/crontab/crontab.pbench
MAILTO="pbench-admins@example.com"
MAILFROM="pbench@example.com"
# 01 0 * * *  /home/pbench/bin/patch_html_cron_script /pbench/public_html/incoming/
# 41 3 * * *  flock -n /home/pbench/locks/copy-sosreport.lock /home/pbench/bin/copy-sosreports
 * * * * *  flock -n /home/pbench/locks/pbench-unpack-tarballs.lock /home/pbench/bin/pbench-unpack-tarballs /pbench
41 * * * *  flock -n /home/pbench/locks/pbench-copy-sosreports.lock /home/pbench/bin/pbench-copy-sosreports /pbench
 * * * * *  flock -n /home/pbench/locks/pbench-edit-prefixes.lock /home/pbench/bin/pbench-edit-prefixes /pbench
# 11 * * * *  flock -n /home/pbench/locks/pbench-create-results-links.lock /home/pbench/bin/pbench-create-results-links /pbench
26 * * * *  flock -n /home/pbench/locks/pbench-index.lock /home/pbench/bin/pbench-index /pbench
51 3 * * *  flock -n /home/pbench/locks/pbench-clean-up-dangling-results-links.lock /home/pbench/bin/pbench-clean-up-dangling-results-links /pbench
